# qr3_46.conf
[program:qr3_46]
command=python3 /home/curator3/qr3_46/main.py
autostart=true
autorestart=true
stderr_logfile=/home/curator3/qr3_46/test.err.log
stdout_logfile=/home/curator3/qr3_46/test.err.log
